# configs/com.fineprint.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fineprint.app</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/fine-print</string>
        <string>-config</string>
        <string>/etc/fine-print/config.yml</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/fine-print/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/fine-print/stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/var/lib/fine-print</string>
    <key>UserName</key>
    <string>root</string>
</dict>
</plist>
